# autostart/lxqt-panel.desktop.in
[Desktop Entry]
Type=Application
TryExec=lwqt-panel
Exec=lwqt-panel
OnlyShowIn=LWQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
